# systemd/investment-weekly.timer
[Unit]
Description=Run Investment Monitor Weekly Synthesis

[Timer]
OnCalendar=Sun 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
