# systemd/pi/ollama.service.d/override.conf
# Drop-in for the official ollama.service — reproduces the macOS RAM guardrails so a
# 16GB Pi never holds two models resident at once (phi3:mini + the 14b synthesis
# model). Installed to /etc/systemd/system/ollama.service.d/override.conf by
# scripts/install.sh, then `systemctl daemon-reload && systemctl restart ollama`.
[Service]
Environment=OLLAMA_MAX_LOADED_MODELS=1
Environment=OLLAMA_NUM_PARALLEL=1
Environment=OLLAMA_KEEP_ALIVE=5m
